# nvidia-uvm-reload/nvidia-uvm-reload.service
[Unit]
Description=Restart nvidia_uvm on dbus method call

[Service]
ExecStart=/usr/bin/nvidia-uvm-reload

[Install]
WantedBy=multi-user.target
