# play.sh/play.sh.service
[Unit]
Description=Watch play.sh config and generate files

[Service]
ExecStart=/usr/bin/play.sh watch

[Install]
WantedBy=default.target